# ops/launchd/com.agentic-trader.premarket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agentic-trader.premarket</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__WRAPPER__</string>
    <string>chain-premarket</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__DEPLOY_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <!-- StartCalendarInterval fires in the MACHINE'S local timezone; install.sh
       verifies the machine is on America/New_York before installing. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>__DEPLOY_ROOT__/logs/launchd/com.agentic-trader.premarket.log</string>
  <key>StandardErrorPath</key>
  <string>__DEPLOY_ROOT__/logs/launchd/com.agentic-trader.premarket.log</string>
</dict>
</plist>
